# docs/systemd/cloudflared-media.service.example
[Unit]
Description=Cloudflare Tunnel for a Plex Live TV proxy frontend
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-networks/
After=network-online.target plex-live-tv-proxy.service
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/cloudflared --config /etc/cloudflared/media-tunnel.yml tunnel run
Restart=always
RestartSec=5s
TimeoutStartSec=90s
KillMode=mixed
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/etc/cloudflared
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
