# docs/systemd/plextuner.service.example
# Example systemd unit for one-run Plex Tuner (Live TV/DVR).
# Copy to /etc/systemd/system/plextuner.service and enable:
#   sudo cp plextuner.service.example /etc/systemd/system/plextuner.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now plextuner
#
# Prerequisites:
#   - .env with PLEX_TUNER_PROVIDER_* (or PLEX_TUNER_M3U_URL) and PLEX_TUNER_BASE_URL
#   - PLEX_TUNER_BASE_URL must be the URL Plex will use (e.g. http://YOUR_SERVER_IP:5004)
#
# One-time in Plex: Settings > Live TV & DVR > Set up > enter Base URL and XMLTV URL (printed at startup).

[Unit]
Description=Plex Tuner (Live TV/DVR)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/plextuner
EnvironmentFile=/opt/plextuner/.env
ExecStart=/opt/plextuner/plex-tuner run -addr :5004 -refresh 6h
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plextuner

# If Base URL is not in .env, set it here (use your server IP or hostname)
# Environment="PLEX_TUNER_BASE_URL=http://192.168.1.10:5004"

[Install]
WantedBy=multi-user.target
